# projects/pigss/build/back-end/lologger/DEBIAN/postinst
#!/bin/bash

sudo mkdir -p /var/log/lologger
sudo chown -R $SUDO_USER:$SUDO_USER /var/log/lologger
sudo systemctl daemon-reload
sudo systemctl enable lologger.service
sudo systemctl start lologger.service
